# etc/services.d/mot_detect/run
#!/command/with-contenv bash
python3 /ext/mot.v.2.1.py $MOT_CAMERANAME $MONGO_USER $MONGO_PASS
