# etc/services.d/mot_detect/run
#!/bin/sh
python3 /ext/ffm.v.2.1.py $MOT_CAMERANAME $MONGO_USER $MONGO_PASS
